# bsc-upgrade.service
[Unit]
Description=Geth Update Service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/bsc_upgrade.sh

[Install]
WantedBy=multi-user.target
