# bsc-upgrade.timer
[Unit]
Description=Timer for Geth Update Service

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
